# package/vwifi/files/etc/uci-defaults/99-vwifi-hwsim
#!/bin/sh
echo "mac80211_hwsim radios=0" > /etc/modules.d/mac80211-hwsim
